# testing/frr/frr-init
#!/sbin/openrc-run

description="Free Range Routing daemon"

start() {
        /usr/lib/frr/frr start
}

stop() {
        /usr/lib/frr/frr stop
}

restart() {
        /usr/lib/frr/frr restart
}

statui() {
        /usr/lib/frr/frr status
}

